# packages/FreeBSD/dns/dnscrypt-proxy/files/dnscrypt-proxy.sh.in
#!/bin/sh

# PROVIDE: dnscrypt-proxy
# REQUIRE: NETWORKING SERVERS
# BEFORE: DAEMON
# KEYWORD: shutdown

#
# Add the following lines to /etc/rc.conf to enable dnscrypt-proxy:
#
# dnscrypt_proxy_enable="YES"
#

. /etc/rc.subr

name=dnscrypt_proxy
rcvar=`set_rcvar`

load_rc_config $name

command=%%PREFIX%%/sbin/dnscrypt-proxy
pidfile=/var/run/dnscrypt-proxy.pid

dnscrypt_proxy_enable=${dnscrypt_proxy_enable:-"NO"}
command_args="${dnscrypt_proxy_config} --daemonize --pidfile=${pidfile}"

run_rc_command "$1"
